# src/daemon/systemd/wifi-telemetryd-ap@.service.in
[Unit]
Description=Wi-Fi Telemetry Daemon (WPA Access Point)
Requires=sys-subsystem-net-devices-%i.device
After=sys-subsystem-net-devices-%i.device

[Service]
Type=simple
Restart=on-failure
ExecStart=${CMAKE_INSTALL_FULL_SBINDIR}/wifi-telemetryd -s wpa -o ap -i %i
