# src/daemon/systemd/wifi-telemetryd-ztp.service.in
[Unit]
Description=Wi-Fi Telemetry Daemon (Zero Touch Provisioning)

[Service]
Type=simple
Restart=on-failure
ExecStart=${CMAKE_INSTALL_FULL_SBINDIR}/wifi-telemetryd -s ztp
